# redux-tutorial/myscript.service
[Unit]
Description=My Startup Script

[Service]
ExecStart=/my-app/redux-tutorial/myscript.sh

[Install]
WantedBy=default.target
